# crontab
# Execute the monitoring script every 30 minutes.

*/30 * * * * /usr/local/bin/php /path/to/monitor.php

